# deploy/research-webhook.service
[Unit]
Description=Research webhook API (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=researcher
Group=researcher
WorkingDirectory=/home/researcher/research-webhook
EnvironmentFile=/home/researcher/research-webhook/.env
# Secrets (WEBHOOK_API_KEY, ANTHROPIC_OAT, ROUTINE_*_FIRE_URL, etc.) live in
# .env. Chmod 600 and owned by researcher. See INFRA.md.
ExecStart=/home/researcher/research-webhook/.venv/bin/uvicorn webhook:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
